# modules/minecraft/files/minecraft.service
[Unit]
Description=minecraft

[Service]
WorkingDirectory=/opt/minecraft

ExecStart=/etc/alternatives/jre_1.8.0_openjdk/bin/java -Xmx1024M -Xms1024M -jar minecraft_server.1.19.3.jar nogui

[Install]
WantedBy=multi-user.target
